# packaging/hearth.openrc
#!/sbin/openrc-run
# OpenRC service for hearth (self-hosted home-automation hub).
#
# Install (as root):
#   sh /datar/workspace/claude_code_experiments/hearth/packaging/install.sh
#
# Watch:  tail -f /datar/workspace/claude_code_experiments/hearth/hearth.log

name="hearth"
description="hearth home-automation hub"

supervisor="supervise-daemon"
command="/datar/workspace/claude_code_experiments/hearth/target/release/hearth"
# Run as the owner of config.toml/token_store.json so token refresh can rewrite them.
command_user="bryan:users"
directory="/datar/workspace/claude_code_experiments/hearth"
output_log="/datar/workspace/claude_code_experiments/hearth/hearth.log"
error_log="/datar/workspace/claude_code_experiments/hearth/hearth.log"
respawn_delay=5
# No respawn_period/respawn_max: a home hub should keep retrying indefinitely on
# crash rather than give up after N tries — and it sidesteps supervise-daemon's
# "respawn-period vs respawn-delay" warning (harmless, but noisy at start).

export RUST_LOG="info"

depend() {
	need net
	after net
}
